# service/discord-wakeonlan.service
[Unit]
Description=Discord Wake-on-LAN Bot
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=INSTALL_DIR_PLACEHOLDER
ExecStart=NODE_PATH_PLACEHOLDER src/main.js
Restart=always
RestartSec=10
StandardOutput=append:INSTALL_DIR_PLACEHOLDER/logs/out.log
StandardError=append:INSTALL_DIR_PLACEHOLDER/logs/err.log
SyslogIdentifier=discord-wakeonlan

# 環境変数ファイル
EnvironmentFile=INSTALL_DIR_PLACEHOLDER/.env

[Install]
WantedBy=multi-user.target