# service/discord-wakeonlan.service
[Unit]
Description=Discord Wake-on-LAN Bot
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=%h/discord-wakeonlan
ExecStart=NODE_PATH_PLACEHOLDER src/main.js
Restart=always
RestartSec=10
StandardOutput=append:%h/discord-wakeonlan/logs/out.log
StandardError=append:%h/discord-wakeonlan/logs/err.log
SyslogIdentifier=discord-wakeonlan

# 環境変数ファイル
EnvironmentFile=%h/discord-wakeonlan/.env

[Install]
WantedBy=default.target